# crontab
0 10 * * * /bin/bash /home/{server_home_DIR}/Projects/e-menu/cron.sh >> /home/{server_home_DIR}/test.log
0 10 * * * echo "Cron works" >> /home/{server_home_DIR}/test.log

